# includes/etc/supervisor/conf.d/nginx.conf
[program:nginx]
startsecs=0
stopwaitsecs=55
command=nginx -g "daemon off;"
autorestart=true
autostart=true
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log
